# systemd/hddfancontrol.service
[Unit]
Description=HDD fan control
Documentation=man:hddfancontrol(1)
After=hddtemp.service

[Service]
EnvironmentFile=/etc/conf.d/hddfancontrol
ExecStart=/usr/bin/hddfancontrol -v $HDDFANCONTROL_LOG_LEVEL daemon $HDDFANCONTROL_DAEMON_ARGS

# Scheduler priority
CPUSchedulingPolicy=rr
CPUSchedulingPriority=49

# Hardening
ProtectSystem=strict
PrivateTmp=true
ProtectHome=true
SystemCallArchitectures=native
MemoryDenyWriteExecute=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
